# worker.service
[Unit]
Description=Plane Worker Daemon
After=network.target

[Service]
Type=simple
User=worker
Group=worker
WorkingDirectory=/opt/plane
ExecStart=/opt/plane/plane
Restart=always
# Ensure secrets are protected
Environment=Plane__secrets__master_key_path=/etc/worker/master.key
Environment=Plane__secrets__worker_secret=YOUR_WORKER_SECRET_HERE

[Install]
WantedBy=multi-user.target
